# gunicorn.service
[Unit]
Description=Gunicorn daemon for miniapp_api Django application
After=network.target postgresql.service

[Service]
User=root
Group=root
WorkingDirectory=/home/miniapp_expert/api-django
Environment="PATH=/home/miniapp_expert/api-django/venv/bin"
ExecStart=/home/miniapp_expert/api-django/venv/bin/gunicorn \
    miniapp_api.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 2 \
    --timeout 120 \
    --error-logfile /tmp/django-error.log \
    --access-logfile /tmp/django-access.log

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

